# services/config/supervisor/conf.d/crond.conf
[program:crond]
command=cron -f
autostart=true
autorestart=true
killasgroup=true
stopasgroup=true