# services/config/supervisor/conf.d/nginx.conf
[program:nginx]
command=nginx
autostart=true
autorestart=true
killasgroup=true
stopasgroup=true
startretries=5
numprocs=1
startsecs=0
process_name=%(program_name)s_%(process_num)02d